# vendor/cookbooks/custombox/templates/default/thumbor.conf
[program:thumbor]
command=thumbor --port=900%(process_num)01d --conf=/etc/thumbor/thumbor.conf
process_name=thumbor900%(process_num)01d
numprocs=1
user=root
;directory=/home/thumbor/
autostart=true
autorestart=true
startretries=3
stopsignal=TERM
;stdout_logfile=/home/thumbor/logs/thumbor800%(process_num)s.stdout.log
;stdout_logfile_maxbytes=1MB
;stdout_logfile_backups=10
;stderr_logfile=/home/thumbor/logs/thumbor800%(process_num)s.stderr.log
;stderr_logfile_maxbytes=1MB
;stderr_logfile_backups=10